# systemd/minecraft-backup.service
[Unit]
Description=Backup script

[Service]
ExecStart=/home/ubuntu/scripts/minecraft-backup.sh
